# supervisord-mysqld.conf
[program:mysqld]
command=/opt/scripts/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true
